# docker/classicpress/config/s6/php-fpm/finish
#!/bin/sh
pkill -TERM "php-fpm${PHP_VERSION}"
