# docker/classicpress/config/s6/php-fpm/run
#!/bin/sh
exec php-fpm${PHP_VERSION} -F
